# docker/php/manager-worker.conf
[program:manager-worker]
command=sh -c "php /var/www/html/console/worker_manager.php"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
numprocs_start=1
redirect_stderr=true
stdout_logfile=/var/www/worker.log
stopwaitsecs=3600